# ubuntu-proxy-manager.service
[Unit]
Description=Ubuntu Proxy Manager service

[Service]
Type=dbus
BusName=com.ubuntu.ProxyManager
ExecStart=/usr/libexec/ubuntu-proxy-manager

[Install]
Alias=dbus-com.ubuntu.ProxyManager
